# testBed/75-myrule.rules
ACTION=="add", ATTRS{idVendor}=="1199", ATTRS{idProduct}=="****", RUN+="/home/xentest/Documents/USB_Auto_Attach/USB_Manager.sh"
